# bot.service
[Unit]
Description=Telegram Bot
After=network.target

[Service]
User=ubuntu
WorkingDirectory=/home/ubuntu/telegram-gpt-bot
ExecStart=/home/ubuntu/telegram-gpt-bot/venv/bin/python3 bot.py
Restart=always
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
